# AutomataFiddle/config/systemd/system/automatafiddle-ssl-renew.service
[Unit]
Description=Attempts to update the Lets Encrypt SSL certificate

[Service]
Type=simple
ExecStart=/var/AutomataFiddle/config/services/automatafiddle-ssh-renew.sh
